# polymaker-autorun.service
[Unit]
Description=Polymarket Maker Autorun Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/trader/polymarket_api/POLYMARKET_MAKER_copytrade_v2/POLYMARKET_MAKER_AUTO
EnvironmentFile=-/home/trader/polymarket_api/POLYMARKET_MAKER_copytrade_v2/.env
Environment="PATH=/root/.pyenv/versions/poly312/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/.pyenv/versions/poly312/bin/python poly_maker_autorun.py --no-repl
Restart=always
RestartSec=10
StandardOutput=append:/home/trader/polymarket_api/POLYMARKET_MAKER_copytrade_v2/POLYMARKET_MAKER_AUTO/autorun_systemd.log
StandardError=append:/home/trader/polymarket_api/POLYMARKET_MAKER_copytrade_v2/POLYMARKET_MAKER_AUTO/autorun_systemd.log

[Install]
WantedBy=multi-user.target
